# rc.local
#!/bin/sh -e
#
# /etc/rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

mount.vboxsf share /mnt/Share vboxsf
mount.vboxsf Projector_switcher_executor /mnt/Projector_switcher_executor vboxsf
mount.vboxsf ESP8266_FOTA /mnt/ESP8266_FOTA vboxsf
exit 0
